# deployments/supervisor/conf.d/paneval-beat.ini
[program:paneval-beat]
directory=/paneval
command=poetry run celery -A paneval beat -l INFO
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1

stopasgroup=false             ; send stop signal to the UNIX process group (default false)
killasgroup=false             ; SIGKILL the UNIX process group (def false)

redirect_stderr=true                        ; redirect proc stderr to stdout (default false)
stdout_logfile=/data/log/app/%(program_name)s/log.log         ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=0                   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=0                    ; # of stdout logfile backups (default 10)
stdout_capture_maxbytes=1MB                 ; number of bytes in 'capturemode' (default 0)
environment=C_FORCE_ROOT=true

stopwaitsecs=60